# daemontools/tls-reverse-proxy/run
#!/bin/bash

ulimit -n 500000
export LANG=ja_JP.UTF-8

. .env

exec 2>&1
exec ./tls-reverse-proxy
